# homeassistant/scripts/Enigma2/crontab
0 1 * * *     sh /etc/CCcam.update
@reboot       sh /etc/CCcam.update